# attributes.list
/etc
rw-rw-r--
root
root

/bin
rw-rw-r--
root
root

/sbin
rw-rw-r--
root
root

/dev
r--r--r--
root
root

/sys
r--r--r--
root
root

/lib
rw-rw-r--
root
root

/var
rw-rw-r--
root
root

/proc
r--r--r--
root
root

/boot
rw-rw-r--
root
root

/etc/release
rw-rw-r--
root
root

/etc/init.d
rw-rw-r--
root
root

/etc/official
rw-rw-r--
root
root

/etc/id
rw-rw-r--
root
root

/etc/hostname
rw-rw-r--
root
root

/etc/fstab
rw-rw-r--
root
root

/bin/shell
rwxrwxr-x
root
root

/bin/ls
rwxrwxr-x
root
root

/bin/pwd
rwxrwxr-x
root
root

/bin/cat
rwxrwxr-x
root
root

/bin/touch
rwxrwxr-x
root
root

/bin/mount
rwxrwxr-x
root
root

/bin/umount
rwxrwxr-x
root
root

/bin/mkdir
rwxrwxr-x
root
root

/bin/dmesg
rwxrwxr-x
root
root

/sbin/init
rwxrwxr-x
root
root

/sbin/banner
rwxrwxr-x
root
root

/lib/system
rw-rw-r--
root
root

/etc/init.d/0-banner
rw-rw-r--
root
root

/etc/init.d/5-users
rw-rw-r--
root
root

/etc/init.d/100-ui
rw-rw-r--
root
root

/lib/system/usrsetup
rwxrwxr--
root
root

/lib/system/initcfg
rwxrwxr--
root
root

/lib/system/users
rw-rw----
root
root

/lib/system/ui
rwxrwxr--
root
root

